# recipes-application/application-example/files/start_application.service
[Unit]
Description=Start dummy application service

[Service]
ExecStart=/usr/bin/app_sample

[Install]
WantedBy=multi-user.target
